# systemd/user/debug-usage-uptime.service
[Unit]
Description=Debug Usage (Uptime)

[Service]
Type=simple
ExecStart=%h/.local/bin/debug-usage uptime 5
