# systemd/user/debug-usage-detailed.service
[Unit]
Description=Debug Usage (Detailed)

[Service]
Type=simple
ExecStart=%h/.local/bin/debug-usage detailed 5
